# debian/openrc.install.in
debian/tmp/bin                /
debian/tmp/sbin               /
debian/tmp/lib/rc             /lib
debian/tmp/etc                /
debian/tmp/usr/share/man/man8 /usr/share/man
debian/inittab /usr/share/openrc
debian/update-rc.d /usr/share/openrc
